# crontab
12 12 * * * cd /Users/shiyiyin/RV/ && /opt/homebrew/bin/python3 main.py